# deb/install/prerm
#!/bin/bash

systemctl stop waggle-rpi-watchdog.service || true
systemctl disable waggle-rpi-watchdog.service || true
